# config/mailhog/mailhog.conf
description "Mailhog"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

exec /usr/bin/env $(which mailhog)
